# service/gunicorn-indi-allsky.service
[Unit]
Description=indi-allsky gunicorn daemon
Requires=%GUNICORN_SERVICE_NAME%.socket
After=network.target

[Service]
Type=notify
#User=%ALLSKY_USER%
#RuntimeDirectory=%GUNICORN_SERVICE_NAME%
WorkingDirectory=%ALLSKY_DIRECTORY%
ExecStart=%ALLSKY_DIRECTORY%/virtualenv/indi-allsky/bin/gunicorn --config %ALLSKY_ETC%/gunicorn.conf.py indi_allsky.wsgi
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true

[Install]
WantedBy=default.target
