# service/indi-allsky.service
[Unit]
Description=Indi AllSky Service
After=network.target indiserver.service
 
[Service]
EnvironmentFile=%ALLSKY_ETC%/indi-allsky.env
WorkingDirectory=%ALLSKY_DIRECTORY%
ExecStart=%ALLSKY_DIRECTORY%/virtualenv/indi-allsky/bin/python3 allsky.py --log syslog run
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID
RestartSec=5
PrivateTmp=true
UMask=0022

[Install]
WantedBy=default.target
