# service/indiserver.service
[Unit]
Description=Indi Server
 
[Service]
User=%INDISERVER_USER%
PAMName=login
ExecStart=%INDI_DRIVER_PATH%/indiserver indi_simulator_telescope %INDI_CCD_DRIVER%
ExecStop=/bin/kill -TERM $MAINPID
RestartSec=5
 
[Install]
WantedBy=multi-user.target
