# service/indi-allsky.timer
[Unit]
Description=Timer for the indi-allsky service

[Timer]
OnBootSec=2min
